# scripts/supervisor-poller.conf
#
# supervisord configuration for meshmap poller, a.k.a `get-map-info.php`
#
[program:poller]
command=/scripts/poller.sh
numprocs=1
autostart=true
autorestart=true
